# assets/somfy.service.tmpl
[Unit]
Description=Somfy Telis 4 remote controller
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User={{SERVICE_USER}}
Group=gpio
ExecStart={{EXEC_START}}
Restart=on-failure
RestartSec=3
Environment=RUST_LOG=info
StateDirectory=somfy
StateDirectoryMode=0700
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
PrivateTmp=true
DevicePolicy=closed
DeviceAllow={{GPIO_CHIP}} rw
DeviceAllow={{SPI_DEVICE}} rw
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
